# docker/dev/app/supervisor/app.conf
[supervisord]
nodaemon=true

[program:app]

[program:nginx]
